# ultraclaude.service
[Unit]
Description=UltraClaude - Multi-session Claude Code Manager
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/rat/ultraclaude
ExecStart=/home/rat/ultraclaude/venv/bin/python3 main.py start --host 127.0.0.1 --port 8420
Restart=always
RestartSec=10
# Only kill the main process on stop/restart — leave tmux sessions alive
# so they persist and reconnect when the server comes back up
KillMode=process
Environment="PATH=/home/rat/.local/bin:/home/rat/.opencode/bin:/home/rat/.nvm/versions/node/v20.19.6/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/usr/lib/wsl/lib:/mnt/c/Program Files/WindowsApps/MicrosoftCorporationII.WindowsSubsystemForLinux_2.5.9.0_x64__8wekyb3d8bbwe:/mnt/c/Program Files/Microsoft SDKs/Azure/CLI2/wbin:/mnt/c/Program Files (x86)/oh-my-posh/bin/:/mnt/c/WINDOWS/system32:/mnt/c/WINDOWS:/mnt/c/WINDOWS/System32/Wbem:/mnt/c/WINDOWS/System32/WindowsPowerShell/v1.0/:/mnt/c/WINDOWS/System32/OpenSSH/:/mnt/c/Program Files (x86)/NVIDIA Corporation/PhysX/Common:/mnt/c/Program Files/dotnet/:/mnt/c/Program Files/Tesseract-OCR/:/mnt/c/ProgramData/chocolatey/bin:/mnt/c/Program Files/nodejs/:/mnt/c/Program Files/CMake/bin:/mnt/c/Program Files (x86)/Windows Kits/10/Windows Performance Toolkit/:/mnt/c/Program Files/Microsoft SQL Server/170/Tools/Binn/:/mnt/c/Program Files/Microsoft SQL Server/Client SDK/ODBC/170/Tools/Binn/:/mnt/c/Program Files/Docker/Docker/resources/bin:/mnt/c/Program Files/PowerShell/7/:/mnt/c/Program Files/Git/cmd:/mnt/c/Users/ratz/AppData/Local/Programs/Python/Launcher/:/mnt/c/Users/ratz/.cargo/bin:/mnt/c/Users/ratz/AppData/Local/Programs/oh-my-posh/bin/:/mnt/c/Users/ratz/AppData/Local/Microsoft/WindowsApps:/mnt/c/Users/ratz/AppData/Local/Programs/cursor/resources/app/bin:/mnt/c/Users/ratz/AppData/Local/Programs/Espanso:/mnt/c/Program Files/mitmproxy/bin:/mnt/c/Users/ratz/AppData/Local/GitHubDesktop/bin:/mnt/c/Users/ratz/AppData/Roaming/npm:/mnt/c/Users/ratz/.dotnet/tools:/mnt/c/Users/ratz/AppData/Local/PowerToys/:/mnt/c/Users/ratz/AppData/Local/Programs/Antigravity/bin:/snap/bin"
Environment="HOME=/home/rat"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ultraclaude

[Install]
WantedBy=default.target
